# crontab
# Run every sunday at 11 PM
00 23 * * 0 root python /PlexEmail/scripts/plexEmail.py --config /config/config.conf;

